# vendor/github.com/ipfs/go-ipfs/misc/launchd/io.ipfs.ipfs-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>KeepAlive</key>
    <true/>
    <key>Label</key>
    <string>io.ipfs.ipfs-daemon</string>
    <key>ProgramArguments</key>
    <array>
      <string>{{IPFS_BIN}}</string>
      <string>daemon</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
      <key>IPFS_PATH</key>
      <string>{{IPFS_PATH}}</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>

