# deployment/noms-server.conf
start on started elastic-network-interfaces
stop on runlevel [S016]

console output
respawn
respawn limit 10 120
exec /home/ec2-user/deployed/server -ldb=/mnt/my-data/noms_11_12